# etc/cron.d/apt-update
0 7 * * * root /usr/bin/apt-get update > /dev/null
